# crontab
*/10 * * * *	www-data if [ "`find /etc/nginx/passwd/transit -mmin +60 -print0`" ]; then echo > /etc/nginx/passwd/transit ; fi
