# crontab
*/10 * * * *	www-data find /etc/nginx/passwd/transit -mmin +60 -print0 | xargs -0 --no-run-if-empty rm
